# examples/dserver-prune-logs.service.example
[Unit]
Description=Delete dserver log files older than 7 days

[Service]
Type=oneshot
User=dserver
Group=dserver
ExecStart=/var/run/dserver/prune_dserver_logs.sh
